# .config/systemd/user/statusline_weather.timer
[Unit]
Description=Update weather status line hourly

[Timer]
OnCalendar=hourly
RemainAfterElapse=no
