# .config/systemd/user/dunst.service.d/override.conf
[Service]
ConditionEnvironment=DISPLAY
Restart=on-failure
RestartSec=5